# deploy/systemd/edenai-db-backup.service
[Unit]
Description=EdenAI PostgreSQL backup
After=network.target

[Service]
Type=oneshot
User=www-data
Group=www-data
EnvironmentFile=/opt/edenai-teacher/backend/.env
ExecStart=/bin/bash /opt/edenai-teacher/deploy/backup_postgres.sh
