# deploy/systemd/edenai-backend.service
[Unit]
Description=EdenAI Teacher FastAPI
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/edenai-teacher/backend
EnvironmentFile=/opt/edenai-teacher/backend/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/edenai-teacher/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/edenai-teacher/backend

[Install]
WantedBy=multi-user.target
